# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:mcp-read]
command=fastmcp run server.py --transport sse --port 8001
directory=/app/mcp-servers/read-mcp
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-sns]
command=fastmcp run server.py --transport sse --port 8002
directory=/app/mcp-servers/sns-resolver
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-token]
command=fastmcp run server.py --transport sse --port 8003
directory=/app/mcp-servers/token-data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-tx]
command=fastmcp run server.py --transport sse --port 8004
directory=/app/mcp-servers/transaction-executor
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-alerts]
command=fastmcp run server.py --transport sse --port 8005
directory=/app/mcp-servers/price-alerts
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-research]
command=fastmcp run server.py --transport sse --port 8007
directory=/app/mcp-servers/market-research
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mcp-elfa-auto]
command=fastmcp run server.py --transport sse --port 8008
directory=/app/mcp-servers/elfa-auto
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:agent]
command=uvicorn main:app --host 0.0.0.0 --port 7777
directory=/app
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:price-monitor]
command=python main.py
directory=/app/services/price_monitor
autostart=false
autorestart=false
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

